# .workspace/supervisord.conf
[supervisord]
nodaemon = true
logfile = /dev/null
logfile_maxbytes = 0
pidfile = /var/run/supervisord.pid

[inet_http_server]
port=9001

[program:symfony-messenger-consume]
command=php /var/www/app/bin/console messenger:consume async --time-limit=3600
user=sipandai
numprocs=2
startsecs=0
autostart=true
autorestart=true
startretries=10
process_name=%(program_name)s_%(process_num)02d